# systemd/bnb-backup.service
[Unit]
Description=bits&bytes Minecraft Backup

[Service]
Type=oneshot
ExecStart=/home/minecraft/server/scripts/backup.sh

